# chroot/var/lib/pacman/local/mesa-22.2.1-1/files
%FILES%
usr/
usr/include/
usr/include/EGL/
usr/include/EGL/eglextchromium.h
usr/include/EGL/eglmesaext.h
usr/include/GL/
usr/include/GL/internal/
usr/include/GL/internal/dri_interface.h
usr/include/GL/osmesa.h
usr/include/d3dadapter/
usr/include/d3dadapter/d3dadapter9.h
usr/include/d3dadapter/drm.h
usr/include/d3dadapter/present.h
usr/include/gbm.h
usr/lib/
usr/lib/bellagio/
usr/lib/bellagio/libomx_mesa.so
usr/lib/d3d/
usr/lib/d3d/d3dadapter9.so
usr/lib/d3d/d3dadapter9.so.1
usr/lib/d3d/d3dadapter9.so.1.0.0
usr/lib/dri/
usr/lib/dri/armada-drm_dri.so
usr/lib/dri/d3d12_dri.so
usr/lib/dri/etnaviv_dri.so
usr/lib/dri/exynos_dri.so
usr/lib/dri/hx8357d_dri.so
usr/lib/dri/ili9225_dri.so
usr/lib/dri/ili9341_dri.so
usr/lib/dri/imx-dcss_dri.so
usr/lib/dri/imx-drm_dri.so
usr/lib/dri/imx-lcdif_dri.so
usr/lib/dri/ingenic-drm_dri.so
usr/lib/dri/kgsl_dri.so
usr/lib/dri/kirin_dri.so
usr/lib/dri/kms_swrast_dri.so
usr/lib/dri/komeda_dri.so
usr/lib/dri/lima_dri.so
usr/lib/dri/mali-dp_dri.so
usr/lib/dri/mcde_dri.so
usr/lib/dri/mediatek_dri.so
usr/lib/dri/meson_dri.so
usr/lib/dri/mi0283qt_dri.so
usr/lib/dri/msm_dri.so
usr/lib/dri/mxsfb-drm_dri.so
usr/lib/dri/nouveau_dri.so
usr/lib/dri/panfrost_dri.so
usr/lib/dri/pl111_dri.so
usr/lib/dri/r300_dri.so
usr/lib/dri/r600_dri.so
usr/lib/dri/radeonsi_dri.so
usr/lib/dri/rcar-du_dri.so
usr/lib/dri/repaper_dri.so
usr/lib/dri/rockchip_dri.so
usr/lib/dri/st7586_dri.so
usr/lib/dri/st7735r_dri.so
usr/lib/dri/stm_dri.so
usr/lib/dri/sun4i-drm_dri.so
usr/lib/dri/swrast_dri.so
usr/lib/dri/v3d_dri.so
usr/lib/dri/vc4_dri.so
usr/lib/dri/virtio_gpu_dri.so
usr/lib/dri/vmwgfx_dri.so
usr/lib/dri/zink_dri.so
usr/lib/libEGL_mesa.so
usr/lib/libEGL_mesa.so.0
usr/lib/libEGL_mesa.so.0.0.0
usr/lib/libGLX_indirect.so.0
usr/lib/libGLX_mesa.so
usr/lib/libGLX_mesa.so.0
usr/lib/libGLX_mesa.so.0.0.0
usr/lib/libOSMesa.so
usr/lib/libOSMesa.so.8
usr/lib/libOSMesa.so.8.0.0
usr/lib/libgbm.so
usr/lib/libgbm.so.1
usr/lib/libgbm.so.1.0.0
usr/lib/libglapi.so
usr/lib/libglapi.so.0
usr/lib/libglapi.so.0.0.0
usr/lib/pkgconfig/
usr/lib/pkgconfig/d3d.pc
usr/lib/pkgconfig/dri.pc
usr/lib/pkgconfig/gbm.pc
usr/lib/pkgconfig/osmesa.pc
usr/share/
usr/share/drirc.d/
usr/share/drirc.d/00-mesa-defaults.conf
usr/share/glvnd/
usr/share/glvnd/egl_vendor.d/
usr/share/glvnd/egl_vendor.d/50_mesa.json
usr/share/licenses/
usr/share/licenses/mesa/
usr/share/licenses/mesa/LICENSE

